# BingWallpaper.cron
0 8 * * * sh /Users/b51/Pictures/BingWallpaper/BingWallpaper.sh
